# couch.service
[Unit]
After=mysql.service

[Service]
User=hjung
Group=adm
ExecStart=/home/hjung/Project/couch/Couch/Couch.sh

[Install]
WantedBy=multi-user.target

