# client-side/couch.service
[Unit]
After=mysql.service

[Service]
ExecStart=/home/ubuntu/Couch/couch.sh

[Install]
WantedBy=multi-user.target

